# .ebextensions/cleanup_fs_logs.cron
# The newline at the end of this file is extremely important.  Cron won't run without it.
*/10 * * * * root /opt/elasticbeanstalk/bin/cleanup_fs_logs.sh
